# initrd/mkosi.skeleton/usr/lib/systemd/system/systemd-networkd.service.d/10-override.conf
[Install]
WantedBy=
WantedBy=sysinit.target
